# bjorn2scan-agent/bjorn2scan-agent.service
[Unit]
Description=Bjorn2Scan Agent - Host-level security scanning agent
Documentation=https://github.com/bvboe/b2s-go
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/var/lib/bjorn2scan/bin/bjorn2scan-agent
Restart=always
RestartSec=5s

# Filesystem access - read-only root with selective write access
ProtectSystem=strict
ReadWritePaths=/var/lib/bjorn2scan
ReadWritePaths=/var/log/bjorn2scan
ReadWritePaths=/tmp
# Docker socket access (optional - agent detects availability at runtime)
# Running as root provides access to /var/run/docker.sock when present
# The '-' prefix makes this non-fatal if the socket doesn't exist
BindReadOnlyPaths=-/var/run/docker.sock

# Security hardening - read-only access to most of the system
NoNewPrivileges=true
PrivateTmp=false
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ProtectHostname=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bjorn2scan-agent

# Resource limits
LimitNOFILE=65536
TasksMax=4096

[Install]
WantedBy=multi-user.target
